# systemd/halloween-projection.service
[Unit]
Description=Halloween Projection Mapper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change the user below if the repository lives under a different account
User=pi
WorkingDirectory=%h/Halloween-2025/Halloween-2025
ExecStart=/usr/bin/env python3 %h/Halloween-2025/Halloween-2025/src/main.py --broker localhost
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
